# gpu-tweaks.service
[Unit]
Description=gpu-tweaks

[Service]
Type=simple
ExecStart=/usr/local/src/gpu-tweaks/on-startup.sh

[Install]
WantedBy=multi-user.target
